# DEBIAN/postinst
#!/bin/sh

# start gitce
/etc/init.d/gitce-watchers start || exit $?

# add as system startup script
update-rc.d gitce-watchers defaults || exit $?
